# packages/adafruit-ap/DEBIAN/postinst
#!/usr/bin/env bash

# backups
mv /etc/dhcp/dhcpd.conf /etc/dhcp/dhcpd.conf.orig
mv /etc/default/isc-dhcp-server /etc/default/isc-dhcp-server.orig
cp /etc/sysctl.conf /etc/sysctl.conf.orig
mv /etc/default/hostapd /etc/default/hostapd.orig
mv /usr/sbin/hostapd /usr/sbin/hostapd.orig
mv /etc/network/interfaces /etc/network/interfaces.orig

cp /etc/adafruit-ap/dhcpd.conf /etc/dhcp/dhcpd.conf
cp /etc/adafruit-ap/isc-dhcp-server /etc/default/isc-dhcp-server
cp /etc/adafruit-ap/hostapd /etc/default/hostapd
cp /etc/adafruit-ap/hostapd-replacement /usr/sbin/hostapd

# take down wlan0
ifdown wlan0

SSID="Pi_AP"
PASS="raspberry"

get_ssid() {

  read -e -p "Enter the WiFi SSID: " -i $SSID SSID

  if [ ${#SSID} -lt 1 ] || [ ${#SSID} -gt 32 ]; then
    echo "WiFi SSID must be 1..32 characters";
    get_ssid;
  fi

}

get_pass() {

  read -e -p "Enter the WiFi Passphrase: " -i $PASS PASS

  if [ ${#PASS} -lt 8 ] || [ ${#PASS} -gt 63 ]; then
    echo "WiFi Passphrase must be 8..63 characters";
    get_pass;
  fi

}

get_ssid;
get_pass;

cat <<EOF > /etc/hostapd/hostapd.conf
interface=wlan0
driver=rtl871xdrv
ssid=$SSID
hw_mode=g
channel=6
macaddr_acl=0
auth_algs=1
ignore_broadcast_ssid=0
wpa=2
wpa_passphrase=$PASS
wpa_key_mgmt=WPA-PSK
wpa_pairwise=TKIP
rsn_pairwise=CCMP
EOF

# set static ip
ifconfig wlan0 192.168.42.1

cat << "EOF" > /etc/network/interfaces
auto lo
auto eth0
 
iface lo inet loopback
iface eth0 inet dhcp
 
allow-hotplug wlan0
iface wlan0 inet static
 address 192.168.42.1
 netmask 255.255.255.0
 up iptables-restore < /etc/iptables.ipv4.nat

EOF

echo "ssid=$SSID" >> /etc/hostapd/hostapd.conf
echo "wpa_passphrase=$PASS" >> /etc/hostapd/hostapd.conf

if ! grep -Fq "ip_forward=1" /etc/sysctl.conf; then
  echo "net.ipv4.ip_forward=1" >> /etc/sysctl.conf
fi
echo "1" > /proc/sys/net/ipv4/ip_forward

# setup iptable forwarding
iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
iptables -A FORWARD -i eth0 -o wlan0 -m state --state RELATED,ESTABLISHED -j ACCEPT
iptables -A FORWARD -i wlan0 -o eth0 -j ACCEPT
iptables-save > /etc/iptables.ipv4.nat
echo "up iptables-restore < /etc/iptables.ipv4.nat" >> /etc/network/interfaces

# start services and enable on boot
service hostapd start 
service isc-dhcp-server start
update-rc.d hostapd enable 
update-rc.d isc-dhcp-server enable
