# packages/adafruit-ap/DEBIAN/postinst
#!/usr/bin/env bash

service isc-dhcp-server stop
service hostapd stop

# disable on boot
update-rc.d -f hostapd remove
update-rc.d -f isc-dhcp-server remove

# backups
mv /etc/dhcp/dhcpd.conf /etc/dhcp/dhcpd.conf.orig
mv /etc/default/isc-dhcp-server /etc/default/isc-dhcp-server.orig
cp /etc/sysctl.conf /etc/sysctl.conf.orig
mv /etc/default/hostapd /etc/default/hostapd.orig
mv /usr/sbin/hostapd /usr/sbin/hostapd.orig
mv /etc/network/interfaces /etc/network/interfaces.orig
mv /etc/hostapd/hostapd.conf /etc/hostapd/hostapd.orig
rm /etc/init.d/adafruit-ap

cp /etc/adafruit-ap/dhcpd.conf /etc/dhcp/dhcpd.conf
cp /etc/adafruit-ap/isc-dhcp-server /etc/default/isc-dhcp-server
cp /etc/adafruit-ap/hostapd /etc/default/hostapd
cp /etc/adafruit-ap/hostapd-replacement /usr/sbin/hostapd
cp /etc/adafruit-ap/interfaces.ap /etc/network/interfaces
cp /etc/adafruit-ap/adafruit-ap /etc/init.d/adafruit-ap

SSID="Pi_AP"
PASS="raspberry"

get_ssid() {

  read -e -p "Enter the WiFi SSID: " -i $SSID SSID

  if [ ${#SSID} -lt 1 ] || [ ${#SSID} -gt 32 ]; then
    echo "WiFi SSID must be 1..32 characters";
    get_ssid;
  fi

}

get_pass() {

  read -e -p "Enter the WiFi Passphrase: " -i $PASS PASS

  if [ ${#PASS} -lt 8 ] || [ ${#PASS} -gt 63 ]; then
    echo "WiFi Passphrase must be 8..63 characters";
    get_pass;
  fi

}

get_ssid;
get_pass;

cat <<EOF > /etc/hostapd/hostapd.conf
interface=wlan0
driver=rtl871xdrv
ssid=$SSID
hw_mode=g
channel=6
macaddr_acl=0
auth_algs=1
ignore_broadcast_ssid=0
wpa=2
wpa_passphrase=$PASS
wpa_key_mgmt=WPA-PSK
wpa_pairwise=TKIP
rsn_pairwise=CCMP
EOF

service adafruit-ap start
