# packages/adafruit-io-gif/etc/init.d/adafruit-io-gif
#!/usr/bin/env bash
### BEGIN INIT INFO
# Provides:          adafruit-io-gif
# Required-Start:    $local_fs $network $named $time $syslog
# Required-Stop:     $local_fs $network $named $time $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop adafruit-io-gif
# Description:       Enables adafruit.io gif display
### END INIT INFO

ada_gif_off () {

  killall -TERM chromium 2>/dev/null;
  killall -TERM matchbox-window-manager 2>/dev/null;
  sleep 2;
  killall -9 chromium 2>/dev/null;
  killall -9 matchbox-window-manager 2>/dev/null;

  forever stop /etc/adafruit-io-gif/node_source/index.js

}

ada_gif_on () {

  xinit aio_gif_start
  forever start /etc/adafruit-io-gif/node_source/index.js

}

case "$1" in
  start) ada_gif_on  ;;
  stop)  ada_gif_off ;;
esac
